# roles/docker_compose_project/templates/reload.timer.j2
[Unit]
Description=Schedules an automatic reload for the {{ project_name }} docker-compose service.

[Timer]
OnCalendar={{ compose_service_autoupdate_interval }}
